# fastapi_be_server/sync.cron
* * * * * root python /app/sync_meilisearch.py >> /var/log/cron.log 2>&1
